# wifibt/bluetooth-ap6335.conf
# Turn on/off the ap6335 bluetooth

description	"Turn on/off the ap6335 bluetooth"

start on started bluetooth
stop on stopped bluetooth

respawn

script
	SERIAL_DEV=/dev/ttyS0
	if [ -r /etc/default/bluetooth-ap6335.conf ]; then
		. /etc/default/bluetooth-ap6335.conf
	fi
	exec /usr/bin/brcm_patchram_plus --patchram /system/vendor/firmware/bcm4339a0.hcd --baudrate 1500000 --enable_lpm --enable_hci --no2bytes ${SERIAL_DEV}
end script

pre-start script
	echo 1 >/sys/class/rfkill/rfkill0/state
end script

post-stop script
	echo 0 >/sys/class/rfkill/rfkill0/state
end script
